# casino-bot.service
[Unit]
Description=Casino Bot Telegram Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/Users/ahmed/Telegram Axis
Environment=PATH=/Users/ahmed/Telegram Axis/.venv/bin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/Users/ahmed/.vscode/extensions/ms-python.debugpy-2025.10.0-darwin-arm64/bundled/scripts/noConfigScripts:/Users/ahmed/Library/Application Support/Code/User/globalStorage/github.copilot-chat/debugCommand
ExecStart=/Users/ahmed/Telegram Axis/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
